# base_app/example_files/yeah_box_agent.service
[Unit]
Description=Multi Pool
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Restart=on-failure
RestartSec=5s
ExecStart=/yeahBox_agent/yeahBox_agent
ExecStop=/bin/kill -TERM $MAINPID
WorkingDirectory=/yeahBox_agent
LimitNOFILE=102400

[Install]
WantedBy=multi-user.target
